# systemd/kyuden-hourly.service
[Unit]
Description=Kyuden hourly collector (user)
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/kyuden-data-collector
EnvironmentFile=%h/kyuden-data-collector/secrets/kyuden.env

Environment=KYUDEN_STATE=%h/kyuden-data-collector/state/storage_state.json
ExecStart=/usr/bin/flock -n %h/kyuden-data-collector/run/hourly.lock \
  %h/kyuden-data-collector/.venv/bin/python %h/kyuden-data-collector/collector.py \
  -m hourly --db %h/kyuden-data-collector/data/kyuden.sqlite
Restart=on-failure
RestartSec=10s
RuntimeMaxSec=300